# com.youtube.subtitleplus.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.youtube.subtitleplus</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/yangping/YouTube-SubtitlePlus/daemon_server.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/yangping/YouTube-SubtitlePlus</string>
    
    <key>StandardOutPath</key>
    <string>/Users/yangping/YouTube-SubtitlePlus/logs/daemon.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/yangping/YouTube-SubtitlePlus/logs/daemon_error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>Nice</key>
    <integer>1</integer>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/yangping</string>
    </dict>
    
    <!-- 可选：设置自动重启间隔 -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>